# server/conf/linkmanager_supervisor.conf
[program:flaskapp]
command = uwsgi --ini /etc/uwsgi/linkmanager.conf

[program:nginx]
command = service nginx restart
